# deploy/com.jumpcworldcup.caffeinate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- macOS LaunchAgent: keeps the Mac awake so the every-10-min live-submit sweeps
     actually FIRE 24/7. Without this, launchd does not wake a sleeping Mac on a
     StartInterval, so the bot would miss any match whose narrow market window
     (markets open only a few hours pre-kickoff and CLOSE at kickoff) falls while
     the machine sleeps.

     `caffeinate -i` asserts "prevent idle system sleep" on BOTH AC and battery, so
     the every-10-min sweeps keep firing even when unplugged (until the battery dies
     — no software can run a laptop past empty).

     Install:
       cp deploy/com.jumpcworldcup.caffeinate.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.jumpcworldcup.caffeinate.plist
     Stop:
       launchctl unload ~/Library/LaunchAgents/com.jumpcworldcup.caffeinate.plist

     LID CLOSED: caffeinate canNOT keep a Mac awake with the lid shut — that needs a
     system setting (run once, persists across reboots until set back to 0):
       sudo pmset -a disablesleep 1      # all power sources, incl. lid closed
       sudo pmset -a disablesleep 0      # undo
     Caveat: lid closed + on battery = heat with no airflow and a fast drain; a cheap
     always-on cloud box is the robust machine-independent option. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jumpcworldcup.caffeinate</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/caffeinate</string>
    <string>-i</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
